# weth-monitor.service
[Unit]
Description=WETH Price Monitor for Arbitrum Mainnet
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/workspace
ExecStart=/usr/bin/python3 -u /workspace/weth_price_monitor.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/weth-monitor.log
StandardError=append:/var/log/weth-monitor-error.log

[Install]
WantedBy=multi-user.target
